# lsmon.crontab
# m h dom mon dow user  command

# Look for and purge old sessions every 30 minutes
*/5 *     * * *     root   /root/lsmon/lsmon.php -s portal,jhbmon -m dl10
1,21,41 *     * * *     root   /root/lsmon/lsmon.php -s portal,jhbmon -m dl100
3 *     * * *     root   /root/lsmon/lsmon.php -s portal,jhbmon -m dl500
* *     * * *     root   /root/lsmon/lsmon.php -m iftraf
